# game-docker/conf/supervisor/conf.d/supervisor-app.ini
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /home/hanfu/rest/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
